# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:xinetd]
command=/usr/sbin/xinetd -dontfork

[program:dss]
user=nobody
command=/opt/dss/bin/dss
